# supervisord/conf.d/userCreate.conf
; supervisord/userCreate.conf
[program:userCreate]
command=php /consumers/users/userCreate.php
autorestart=true
