# supervisord/supervisord.conf
[supervisord]
nodaemon = False
logfile = /var/log/supervisor/supervisord.log
pidfile = /var/run/supervisord.pid
childlogdir = /var/log/supervisor

[program:userCreate]
command=php /var/www/html/consumers/users/userCreate.php
autorestart=true

[program:userDelete]
command=php /var/www/html/consumers/users/userDelete.php
autorestart=true

[program:userUpdate]
command=php /var/www/html/consumers/users/userUpdate.php
autorestart=true


[program:ucCreate]
command=php /var/www/html/consumers/ucs/ucCreate.php
autorestart=true

[program:ucDelete]
command=php /var/www/html/consumers/ucs/ucDelete.php
autorestart=true

[program:ucUpdate]
command=php /var/www/html/consumers/ucs/ucUpdate.php
autorestart=true